# email/docker/mail/supervisord.conf
; Supervisor config for the FlowOne mail pod. Runs the full native mail stack as
; foreground, auto-restarting programs. Each program's stdio is forwarded to the
; container stdout/stderr so `docker logs flowone-mail` shows everything.
;
; Core mail flow (postfix/dovecot/opendkim) does not hard-depend on the security
; services: Postfix is milter_default_action=accept, so mail still flows while
; clamav/rspamd/spamassassin warm up or restart.

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:rsyslog]
command=/usr/sbin/rsyslogd -n -f /etc/rsyslog.conf
priority=1
autorestart=true
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:unbound]
command=/usr/sbin/unbound -d -p
priority=5
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:redis-rspamd]
command=/usr/bin/redis-server /etc/redis/rspamd-redis.conf
priority=5
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:clamd]
command=/usr/sbin/clamd --foreground
priority=10
; clamd refuses to start until freshclam has fetched a virus DB; let it retry.
autorestart=true
startsecs=5
startretries=1000
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:freshclam]
command=/bin/sh -c "while :; do /usr/bin/freshclam --stdout --no-warnings || true; sleep 3600; done"
priority=10
autorestart=true
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:opendkim]
command=/usr/sbin/opendkim -f -x /etc/opendkim.conf
priority=15
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:opendmarc]
command=/usr/sbin/opendmarc -f -c /etc/opendmarc.conf
priority=15
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:spamd]
command=/usr/sbin/spamd --nouser-config --helper-home-dir=/var/lib/spamassassin --max-children=5 --syslog=stderr
priority=15
autorestart=true
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:spamass-milter]
command=/usr/local/bin/run-spamass.sh
priority=20
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:rspamd]
command=/usr/bin/rspamd -f -u _rspamd -g _rspamd
priority=20
autorestart=true
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:dovecot]
command=/usr/sbin/dovecot -F
priority=25
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:postfix]
command=/usr/sbin/postfix start-fg
priority=30
autorestart=true
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
